# crontab
# comment
*/1 * * * * * echo "hello world!"
*/1 * * * * * echo "hello woensdag" ; sleep 1 ; echo "world"
